# deploy/ec2/systemd/blog-scheduler.service
[Unit]
Description=Jaubi Blog Scheduler (Laravel schedule:run)
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/jaubi.co.kr/blog/backend
ExecStart=/usr/bin/php artisan schedule:run
